# scripts/vulndesk.service
[Unit]
Description=Vulndesk
Documentation=https://github.com/SauceTaster/Vulndesk
After=network.target

[Service]
Environment=NODE_ENV=production
Type=simple
User=webapp
WorkingDirectory=/data/server/repos/Vulndesk
# The app is TypeScript; run it through tsx (installed as a dependency).
ExecStart=/data/server/repos/Vulndesk/node_modules/.bin/tsx app.ts
Restart=on-failure

[Install]
WantedBy=multi-user.target
